# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:mysql]
command=/usr/sbin/mysqld
autostart=true
autorestart=true

[program:php]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
